# overlay/system/etc/systemd/system/usb-moded.service.d/debug-log.conf
[Service]
StandardOutput=append:/userdata/usb-moded.log
StandardError=append:/userdata/usb-moded.log
