# overlay/system/etc/systemd/system/usb-moded.service.d/gadget-reset.conf
[Service]
ExecStartPre=-/bin/sh /userdata/usb-gadget-preconfigure.sh
