# overlay/system/etc/systemd/system/usb-moded.service.d/zz-usb-gadget-reset.conf
[Service]
Type=simple
ExecStartPre=
ExecStartPre=/bin/sh /userdata/usb-moded-conf-wrapper.sh
Environment=RNDIS_FUNCTION=rndis.usb0
Environment=RNDIS_IFNAME=rndis0
Environment=USB_MODED_ARGS=

ExecStart=
ExecStart=/bin/sh /userdata/usb-moded-bypass.sh

Restart=no

# usb_moded's own CapabilityBoundingSet (see rootfs usb-moded.service)
# never included CAP_SYS_CHROOT/CAP_CHOWN/CAP_FOWNER since sshd was never
# meant to run under this unit. We now run sshd from here:
#  - privsep preauth does chroot("/run/sshd") unconditionally -> needs
#    CAP_SYS_CHROOT, else "chroot(...): Operation not permitted" and the
#    connection gets reset.
#  - interactive sessions (pty allocation) chown(/dev/pts/N) to the login
#    user -> needs CAP_CHOWN (and CAP_FOWNER), else that fails and the
#    session gets closed right after auth. Non-pty command exec doesn't
#    hit this path, which is why single-command ssh calls worked already.
#  - login session accounting writes a netlink audit record -> needs
#    CAP_AUDIT_WRITE, else "linux_audit_write_entry failed" is fatal.
# This unit now hosts our debug SSH shell, which keeps hitting new
# missing-capability walls one at a time (chroot, chown, audit, ptrace,
# setpcap...) as we use it for actual system debugging, not just USB
# gadget setup. Instead of chasing them individually, grant the full
# capability set so root-via-SSH behaves like a normal root shell.
# NOTE: this systemd version doesn't merge CapabilityBoundingSet lines,
# so the whole set must be listed together on one/consecutive lines.
CapabilityBoundingSet=
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER CAP_FSETID CAP_KILL CAP_SETGID CAP_SETUID CAP_SETPCAP CAP_LINUX_IMMUTABLE CAP_NET_BIND_SERVICE CAP_NET_BROADCAST CAP_NET_ADMIN CAP_NET_RAW CAP_IPC_LOCK CAP_IPC_OWNER CAP_SYS_MODULE CAP_SYS_RAWIO CAP_SYS_CHROOT CAP_SYS_PTRACE CAP_SYS_PACCT CAP_SYS_ADMIN CAP_SYS_BOOT CAP_SYS_NICE CAP_SYS_RESOURCE CAP_SYS_TIME CAP_SYS_TTY_CONFIG CAP_MKNOD CAP_LEASE CAP_AUDIT_WRITE CAP_AUDIT_CONTROL CAP_SETFCAP CAP_MAC_OVERRIDE CAP_MAC_ADMIN CAP_SYSLOG CAP_WAKE_ALARM CAP_BLOCK_SUSPEND CAP_AUDIT_READ
