# overlay/system/etc/systemd/system/fix-system-vendor-symlinks.service
[Unit]
Description=Fix host-side /system and /vendor symlinks (double-nested system.img layout)
After=lxc-android-config.service
Requires=lxc-android-config.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh /userdata/fix-system-vendor-symlinks.sh

[Install]
WantedBy=multi-user.target
