# overlay/system/etc/systemd/system/lxc-android-config.service.d/cleanup.conf
[Service]
ExecStartPre=/bin/sh /userdata/lxc-mount-cleanup.sh
